# debian/postinst
#!/bin/sh
# postinst script for Rivolution
set -e

case "$1" in
    configure)
      # $2 is the previously-configured version string -- empty exactly
      # when this is a fresh install (Debian policy's own convention for
      # postinst's "configure" argument). Used below to gate destructive
      # database operations to first-install only.
      OLD_VERSION="$2"

      ldconfig

      # Load-bearing, easy to skip silently: ldconfig processes
      # /etc/ld.so.conf.d/*.conf alphabetically, and the standard
      # multiarch conf sorts before pipewire-jack's -- so a real
      # libjack-jackd2-0, if present, always wins SONAME resolution and
      # PipeWire audio silently never connects despite every service
      # reporting active.
      #
      # pipewire-jack ships this conf file only as a doc/example
      # (/usr/share/doc/pipewire/examples/ld.so.conf.d/), never installed
      # to /etc/ld.so.conf.d/ automatically -- confirmed 2026-07-02 on a
      # real fresh install, where the rename-only version of this step
      # silently found nothing to rename (the source file never existed
      # in /etc/ld.so.conf.d/ in the first place) and JACK never worked.
      # Copy it into place with the 00- prefix directly, then also handle
      # the case where an older deployment already copied it under its
      # original (unprefixed) name.
      #
      # Multiarch triplet read straight off the filesystem rather than
      # via `dpkg-architecture` -- that command lives in dpkg-dev, a
      # build-time-only package. It's present on any box that can build
      # this .deb (which is why this passed local testing) but not on a
      # plain install target: confirmed 2026-07-02 via a real install,
      # "dpkg-architecture: not found", postinst exit 127.
      # /usr/lib/<triplet>/ always exists on any multiarch Debian/Ubuntu
      # system regardless of what's installed.
      multiarch_dir=$(ls -d /usr/lib/*-linux-gnu 2>/dev/null | head -n1)
      if test -n "$multiarch_dir" ; then
	  pw_jack_conf="pipewire-jack-$(basename "$multiarch_dir").conf"
	  if test -e "/usr/share/doc/pipewire/examples/ld.so.conf.d/$pw_jack_conf" && \
	     test ! -e "/etc/ld.so.conf.d/00-$pw_jack_conf" ; then
	      cp "/usr/share/doc/pipewire/examples/ld.so.conf.d/$pw_jack_conf" "/etc/ld.so.conf.d/00-$pw_jack_conf"
	  fi
      fi
      for f in /etc/ld.so.conf.d/pipewire-jack-*.conf; do
	  test -e "$f" || continue
	  target="$(dirname "$f")/00-$(basename "$f")"
	  test -e "$target" || mv "$f" "$target"
      done
      ldconfig

      /bin/systemctl daemon-reload

      # rd: the service account this fork's spec 0007/0010 model runs
      # caed/liquidsoap/rivapi/stereo-tool as. Expected to already exist
      # (the operator's own login on a real install) -- this package does
      # not create human accounts, only grants the audio group needed for
      # real-time sound device access if the account is present.
      if getent passwd rd >/dev/null; then
	  usermod -aG audio rd

	  # Stereo Tool reaches JACK only via ALSA's "jack (ALSA)" I/O
	  # option, which routes through the stock pcm.jack definition
	  # (shipped by libasound2-plugins, a Depends) -- that stock
	  # definition hardcodes system:capture_1/2 / system:playback_1/2,
	  # real jackd+hardware ALSA port names that don't exist under
	  # this fork's system-scope PipeWire, causing Stereo Tool's ALSA
	  # layer to fail outright rather than become a valid, patchable
	  # JACK client. This override points it at real, existing ports
	  # instead. /patchbay's reconciler can still change the *actual*
	  # routing afterward -- this only gives Stereo Tool something
	  # valid to open at all. Confirmed 2026-07-02 missing on a real
	  # install: not deploying this was a mistake in the original
	  # postinst automation, corrected here.
	  rd_home=$(getent passwd rd | cut -d: -f6)
	  if test -n "$rd_home" && test -d "$rd_home" && test ! -e "$rd_home/.asoundrc" ; then
	      cp /usr/share/rivolution/alsa/rd.asoundrc "$rd_home/.asoundrc"
	      chown rd:rd "$rd_home/.asoundrc" 2>/dev/null || true
	  fi
      fi

      # Legacy rivendell/pypad system accounts -- upstream's fixed
      # UID/GID 150/151, still referenced by file ownership below
      # (rdselect_helper, webget.cgi, /var/log/rivendell) even though the
      # actual broadcast services now run as rd.
      /usr/sbin/groupadd -r -g 150 rivendell
      /usr/sbin/useradd -o -u 150 -g rivendell -s /bin/false -r -c "Rivendell radio automation system" -d /var/snd rivendell
      /usr/sbin/groupadd -r -g 151 pypad
      /usr/sbin/useradd -o -u 151 -g pypad -s /bin/false -r -c "Rivendell PyPAD scripts" -d /dev/null pypad

      # /var/snd was owned rivendell:rivendell 775 -- rd (the account
      # caed/rdimport/rdlibrary actually run as, per the spec 0007/0010
      # User=rd migration) was never a member of that group, so it only
      # got the "other" bits there (r-x, no write). Audio import failed
      # outright ("Unable to create destination file") with no
      # obviously-wrong permission bit to spot by eye -- confirmed
      # 2026-07-02 on a real install. Fixed to rd:rivendell -- rd owns
      # it outright (native rdimport/caed writes work regardless of
      # group), rivendell keeps write via the group (rdxport.cgi's
      # setuid-drop identity), and unlike relying on supplementary group
      # membership alone, this also lets a third-party sync client
      # (e.g. a desktop sync tool staging files in as plain rd) write
      # into /var/snd without needing group membership synced too.
      # Applied unconditionally (not just on first creation) so an
      # upgrade from an older rivendell:rivendell install self-heals.
      if getent passwd rd >/dev/null; then
	  usermod -aG rivendell rd
      fi

      mkdir -p /var/snd
      chown rd:rivendell /var/snd
      chmod 775 /var/snd
      if test ! -d /etc/rivendell.d ; then
	  mkdir -p /etc/rivendell.d
	  chmod 775 /etc/rivendell.d
      fi
      if test ! -e /etc/rd.conf ; then
	  cp /usr/share/rivolution/rd.conf-sample /etc/rivendell.d/rd-default.conf
	  ln -s /etc/rivendell.d/rd-default.conf /etc/rd.conf
      fi
      if test ! -h /etc/rd.conf ; then
	  mv /etc/rd.conf /etc/rivendell.d/rd-default.conf
	  ln -s /etc/rivendell.d/rd-default.conf /etc/rd.conf
      fi
      if test ! -e /etc/profile.d/rivendell-env.sh ; then
	 cp /usr/share/rivolution/rivendell-env.sh /etc/profile.d/
      fi

      #
      # FIXME: Configure ALSA Here!
      #

      # Database: a fresh install creates and seeds an empty database
      # (mirrors scripts/rivolution-first-run.sh, folded in here so a
      # .deb install needs no separate manual step); an upgrade only
      # migrates the schema against the existing database. This must run
      # -- and on upgrade, complete -- before rivendell.service is
      # restarted below: a schema/binary version mismatch hard-blocks
      # startup (rdcoreapplication.cpp:240).
      if test -z "$OLD_VERSION" ; then
	  rd_conf_real=/etc/rivendell.d/rd-default.conf

	  mysql_pass=$(tr -dc 'A-Za-z0-9' < /dev/urandom | head -c 32)
	  sed -i "/^\[mySQL\]/,/^\[/{s/^Password=.*/Password=$mysql_pass/}" "$rd_conf_real"

	  jwt_secret=$(tr -dc 'A-Za-z0-9' < /dev/urandom | head -c 64)
	  sed -i "/^\[dashboard\]/,/^\[/{s/^JwtSecret=.*/JwtSecret=$jwt_secret/}" "$rd_conf_real"

	  mysql_host=$(awk -F= '/^\[mySQL\]/{f=1;next}/^\[/{f=0}f&&$1=="Hostname"{print $2}' /etc/rd.conf)
	  mysql_user=$(awk -F= '/^\[mySQL\]/{f=1;next}/^\[/{f=0}f&&$1=="Loginname"{print $2}' /etc/rd.conf)
	  mysql_db=$(awk -F= '/^\[mySQL\]/{f=1;next}/^\[/{f=0}f&&$1=="Database"{print $2}' /etc/rd.conf)

	  mysql -h "$mysql_host" <<SQL
drop database if exists \`$mysql_db\`;
create database if not exists \`$mysql_db\`;
drop user if exists '$mysql_user'@'%';
drop user if exists '$mysql_user'@'localhost';
create user '$mysql_user'@'%' identified by '$mysql_pass';
create user '$mysql_user'@'localhost' identified by '$mysql_pass';
grant select,insert,update,delete,create,drop,index,alter,lock tables on \`$mysql_db\`.* to '$mysql_user'@'%';
grant select,insert,update,delete,create,drop,index,alter,lock tables on \`$mysql_db\`.* to '$mysql_user'@'localhost';
flush privileges;
SQL
	  /usr/sbin/rddbmgr --create --generate-audio
      else
	  /usr/sbin/rddbmgr --modify
      fi

      # PulseAudio disabled in favor of direct ALSA/PipeWire access --
      # caed and the system-scope PipeWire instance need uncontended
      # access to the sound device.
      killall pulseaudio >/dev/null 2>&1 || true
      if test -e /etc/pulse/client.conf ; then
	  sed -i 's/# autospawn = yes/autospawn = no/' /etc/pulse/client.conf
      fi
      gpasswd -d pulse audio >/dev/null 2>&1 || true
      usermod -aG audio rivendell
      if ! grep -qF '@audio      hard      memlock     unlimited' /etc/security/limits.conf ; then
	  cat >> /etc/security/limits.conf <<'LIMITS'
@audio      hard      rtprio          90
@audio      hard      memlock     unlimited
LIMITS
      fi

      #
      # Broadcast/PipeWire/rivapi runtime layer (specs 0007/0008/0010).
      # Staged files copied from /usr/share/rivolution/ into their real
      # /etc destinations, same pattern as rd.conf-sample above.
      #
      mkdir -p /etc/systemd/system/rivendell.service.d
      mkdir -p /etc/systemd/system/icecast2.service.d
      mkdir -p /etc/systemd/system/liquidsoap.service.d
      cp /usr/share/rivolution/systemd/rivolution-stack.target /etc/systemd/system/
      cp /usr/share/rivolution/systemd/pipewire-system.service /etc/systemd/system/
      cp /usr/share/rivolution/systemd/wireplumber-system.service /etc/systemd/system/
      cp /usr/share/rivolution/systemd/liquidsoap.service /etc/systemd/system/
      cp /usr/share/rivolution/systemd/stereo-tool.service /etc/systemd/system/
      cp /usr/share/rivolution/systemd/rivapi.service /etc/systemd/system/
      cp /usr/share/rivolution/systemd/rivendell.service.d/rivolution.conf /etc/systemd/system/rivendell.service.d/
      cp /usr/share/rivolution/systemd/icecast2.service.d/rivolution.conf /etc/systemd/system/icecast2.service.d/
      cp /usr/share/rivolution/systemd/liquidsoap.service.d/rivolution.conf /etc/systemd/system/liquidsoap.service.d/

      cp /usr/share/rivolution/sudoers.d/rivapi /etc/sudoers.d/rivapi
      chmod 440 /etc/sudoers.d/rivapi

      cp /usr/share/rivolution/udev/99-ptp.rules /etc/udev/rules.d/
      udevadm control --reload-rules
      udevadm trigger

      systemctl daemon-reload

      # PipeWire must be up before rivendell.service (caed's JACK driver
      # needs the socket present at startup, not just eventually), so
      # these two start ahead of the rivendell restart below.
      systemctl enable --now pipewire-system.service wireplumber-system.service
      systemctl enable --now rivapi.service
      systemctl enable rivolution-stack.target

      /bin/systemctl restart rivendell
      /bin/systemctl enable rivendell
      ln -f -s ../conf-available/rd-bin.conf /etc/apache2/conf-enabled/rd-bin.conf
      ln -f -s ../mods-available/cgid.conf /etc/apache2/mods-enabled/cgid.conf
      ln -f -s ../mods-available/cgid.load /etc/apache2/mods-enabled/cgid.load
      chown rivendell:rivendell /var/www/rd-bin/rdxport.cgi
      chmod 6755 /var/www/rd-bin/rdxport.cgi
      /bin/systemctl restart apache2
      /bin/systemctl enable apache2

      if test -x /usr/bin/gtk-update-icon-cache ; then
	  /usr/bin/gtk-update-icon-cache -f --quiet /usr/share/icons/hicolor
      fi
      if test ! -e /var/www/rd-bin/logos ; then
	  mkdir -p /var/www/rd-bin/logos
      fi
      if test ! -f /var/www/rd-bin/logos/webget_logo.png ; then
	  mkdir -p /var/www/rd-bin/logos
	  cp -f /usr/share/rivolution/logos/webget_logo.png /var/www/rd-bin/logos/webget_logo.png
      fi

      mkdir -p /var/log/rivendell
      chmod 750 /var/log/rivendell
      chown root:syslog /var/log/rivendell
      touch /var/log/rivendell/operations
      chmod 640 /var/log/rivendell/operations
      chown syslog:adm /var/log/rivendell/operations
      touch /var/log/rivendell/debug
      chmod 640 /var/log/rivendell/debug
      chown syslog:adm /var/log/rivendell/debug
      if test ! -e /etc/rsyslog.d/10-rivendell.conf ; then
	  cp /usr/share/rivolution/syslog.conf-sample /etc/rsyslog.d/10-rivendell.conf
      fi
      /bin/systemctl restart rsyslog

      if test -x /usr/bin/rdselect_helper ; then
	  chown root:root /usr/bin/rdselect_helper
	  chmod 4755 /usr/bin/rdselect_helper
      fi

      if test -x /var/www/rd-bin/webget.cgi ; then
	 chown root:root /var/www/rd-bin/webget.cgi
	 chmod 6755 /var/www/rd-bin/webget.cgi
      fi

      echo "Rivolution installed. One manual step remains -- it's a"
      echo "browser action, not a terminal command: open the dashboard"
      echo "(http://<this-host>:8080/), go to /patchbay, connect and save"
      echo "the caed -> Stereo Tool -> Liquidsoap audio chain."
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
