# conf/systemd/icecast2.service.d/rivolution.conf
# Drop-in for icecast2.service: start after Rivendell is fully ready
# (caed must be accepting connections before Liquidsoap can source audio,
# and Icecast must be up before Liquidsoap tries to push a stream to it).
#
# Install: copy to /etc/systemd/system/icecast2.service.d/rivolution.conf
# then: sudo systemctl daemon-reload

[Unit]
After=rivendell.service
PartOf=rivolution-stack.target
