# conf/systemd/liquidsoap.service.d/rivolution.conf
# Drop-in for liquidsoap.service: start after both Rivendell and Icecast
# are fully ready. Liquidsoap sources audio from Rivendell's audio ports
# and pushes the encoded stream to Icecast's source port — both must be
# accepting connections before Liquidsoap starts.
#
# Install: copy to /etc/systemd/system/liquidsoap.service.d/rivolution.conf
# then: sudo systemctl daemon-reload

[Unit]
After=rivendell.service
After=icecast2.service
PartOf=rivolution-stack.target

[Service]
# Clear the upstream ExecStart and point it at the rivapi-generated script.
ExecStart=
ExecStart=/usr/bin/liquidsoap /home/rd/etc/liquidsoap/radio.liq
