# conf/systemd/liquidsoap.service.d/rivolution.conf
# Drop-in for liquidsoap.service.
# The Ubuntu liquidsoap package does not ship a service unit, so this
# drop-in applies on top of conf/systemd/liquidsoap.service (our own unit).
# The After= and PartOf= here are redundant with the main unit but kept
# as a clear statement that this unit belongs to the rivolution stack.
#
# Install: copy to /etc/systemd/system/liquidsoap.service.d/rivolution.conf
# then: sudo systemctl daemon-reload

[Unit]
After=rivendell.service
After=icecast2.service
PartOf=rivolution-stack.target
