# conf/systemd/rivendell.service.d/rivolution.conf
# Drop-in override for rivendell.service (Phase 1).
# Runs rdservice and all its children (caed, ripcd, rdcatchd, etc.) as the
# rd user instead of root. This is the prerequisite for PipeWire integration:
# a root-owned process cannot reach the system-scope PipeWire socket.
#
# Phase 1.5 will split caed into its own standalone caed.service unit with
# per-process resource limits. At that point this drop-in will be revised
# to remove the caed readiness probe (handled by caed.service itself) and
# update After= to reflect the new chain.
#
# Install: copy to /etc/systemd/system/rivendell.service.d/rivolution.conf
# then: sudo systemctl daemon-reload && sudo systemctl restart rivendell.service

[Unit]
PartOf=rivolution-stack.target
# mariadb.service is Type=notify -- After= alone is sufficient to wait
# for it to actually be ready to accept connections, not just started.
# Without this, rivendell.service's first attempt on a fresh boot can
# fail before mariadb is ready; Restart=always then retries and
# succeeds a few seconds later, but any *other* unit relying on bare
# After=rivendell.service (without its own readiness check) can start
# during that failed-then-retrying window, since a failed unit still
# counts as "reached" for ordering purposes (see stereo-tool.service's
# own comment for the case this was first found in).
After=mariadb.service

[Service]
User=rd
Group=rd
# Route caed's JACK calls to the system-scope PipeWire instance.
# pipewire-jack intercepts libjack.so; XDG_RUNTIME_DIR tells it where
# the PipeWire socket lives (/run/pipewire-system/pipewire-0).
Environment=XDG_RUNTIME_DIR=/run/pipewire-system
# Real-time scheduling priority for caed (audio clock). All rdservice
# children inherit this limit; non-audio children don't use SCHED_FIFO
# so the limit is harmless to them.
LimitRTPRIO=99
LimitRTTIME=infinity
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Readiness probe: poll caed's TCP port (5005) until it accepts a
# connection, signalling that caed is fully initialized and the stack
# is ready for downstream services. Polls every second for up to 30s.
# Replace with native sd_notify support in caed when Phase 1.5 lands.
ExecStartPost=/bin/sh -c 'i=0; while [ $i -lt 30 ]; do nc -z 127.0.0.1 5005 2>/dev/null && exit 0; sleep 1; i=$((i+1)); done; echo "caed (port 5005) not ready after 30s" >&2; exit 1'
